# kaonic-factory/kaonic-factory.service
[Unit]
Description=Kaonic Factory
After=network.target kanoic-commd.service

[Service]
ExecStart=/usr/bin/kaonic-factory
Restart=always
RestartSec=5
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

